# deploy/shelf.service
# EXAMPLE unit — illustrative only. `./install.sh` GENERATES the real unit from your actual paths
# with User=$RUN_USER and the venv it created; the /opt/WEB paths below are placeholders. If you
# copy this by hand, edit Documentation=/WorkingDirectory=/ExecStart= to your install dir and set a
# non-root User=. Prefer running install.sh.
[Unit]
Description=Shelf — self-hosted fiction reader (API + UI)
Documentation=file:/opt/WEB/README.md
After=network-online.target
Wants=network-online.target

[Service]
# Type=notify: the app sends READY=1 once it's serving, and pings the watchdog below.
Type=notify
NotifyAccess=main
WatchdogSec=120
WorkingDirectory=/opt/WEB/backend
# Bind on all interfaces; serve the built SPA from ../frontend/dist.
Environment=SHELF_HOST=0.0.0.0
Environment=SHELF_PORT=8000
# Playwright browsers (for per-source render_js) live in root's cache.
Environment=PLAYWRIGHT_BROWSERS_PATH=/root/.cache/ms-playwright
ExecStart=/opt/WEB/backend/.venv/bin/python -m app
Restart=always
RestartSec=3
TimeoutStopSec=20
# Mixed: SIGTERM the main process (graceful), SIGKILL any lingering children.
KillMode=mixed
# Light hardening (relaxed enough for SQLite writes in /opt/WEB and the playwright cache).
NoNewPrivileges=true
ProtectSystem=full
ReadWritePaths=/opt/WEB /root/.cache/ms-playwright

[Install]
WantedBy=multi-user.target
